# debian/gecko-controller/DEBIAN/postinst
#!/bin/sh
set -e

# Source debconf library
. /usr/share/debconf/confmodule

# Function to log messages
log() {
    echo "$1" >&2
    logger -t gecko-controller "$1"
}

configure_i2c() {
    CONFIG="/boot/firmware/config.txt"
    if [ -f "$CONFIG" ]; then
        # Remove any existing i2c_arm settings
        sed -i '/^dtparam=i2c_arm/d' "$CONFIG"

        # Add new i2c configuration with lower baud rate
        echo "# Modified by gecko-controller for long I2C cables" >> "$CONFIG"
        echo "dtparam=i2c_arm=on,i2c_arm_baudrate=10000" >> "$CONFIG"

        log "I2C configuration updated in /boot/firmware/config.txt"
    else
        log "Warning: /boot/firmware/config.txt not found - manual I2C configuration may be needed"
    fi
}

configure_system() {
    # Create gecko-controller user if it doesn't exist
    if ! getent passwd gecko-controller > /dev/null; then
        adduser --quiet --system --home /var/lib/gecko-controller --group gecko-controller
    fi

    # Add gecko-controller user to required groups
    usermod -a -G gpio,i2c,dialout gecko-controller || true

    # Create required directories with proper permissions
    install -d -m 755 -o gecko-controller -g gecko-controller /var/lib/gecko-controller
    install -d -m 755 -o gecko-controller -g gecko-controller /var/log/gecko-controller

    # Ensure config directory permissions are correct
    chown -R gecko-controller:gecko-controller /etc/gecko-controller
    chmod 750 /etc/gecko-controller
    chmod 640 /etc/gecko-controller/config.py
}

case "$1" in
    configure)
        log "Starting gecko-controller post-installation configuration"

        # Configure system user and directories
        configure_system

        # Configure I2C
        configure_i2c

        # Reload systemd in case the service files changed
        if [ -x /bin/systemctl ]; then
            systemctl daemon-reload || true

            # Enable and start services
            systemctl enable gecko-controller.service || true
            systemctl restart gecko-controller.service || true
        fi

        log "Post-installation configuration completed"
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.
# Automatically added by dh_systemd_enable/13.11.4
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	# The following line should be removed in trixie or trixie+1
	deb-systemd-helper unmask 'gecko-controller.service' >/dev/null || true

	# was-enabled defaults to true, so new installations run enable.
	if deb-systemd-helper --quiet was-enabled 'gecko-controller.service'; then
		# Enables the unit on first installation, creates new
		# symlinks on upgrades if the unit file has changed.
		deb-systemd-helper enable 'gecko-controller.service' >/dev/null || true
	else
		# Update the statefile to add new symlinks (if any), which need to be
		# cleaned up on purge. Also remove old symlinks.
		deb-systemd-helper update-state 'gecko-controller.service' >/dev/null || true
	fi
fi
# End automatically added section

# Automatically added by dh_python3
if command -v py3compile >/dev/null 2>&1; then
	py3compile -p gecko-controller 
fi
if command -v pypy3compile >/dev/null 2>&1; then
	pypy3compile -p gecko-controller  || true
fi

# End automatically added section
# Automatically added by dh_systemd_start/13.11.4
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null || true
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		deb-systemd-invoke $_dh_action 'gecko-controller.service' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installinit/13.11.4
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -z "${DPKG_ROOT:-}" ] && [ -x "/etc/init.d/gecko-controller" ]; then
		update-rc.d gecko-controller defaults >/dev/null
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		invoke-rc.d gecko-controller $_dh_action || exit 1
	fi
fi
# End automatically added section


exit 0
