# debian/gecko-controller.service
[Unit]
Description=Gecko Enclosure Controller
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/gecko-controller
Environment=PYTHONPATH=/usr/lib/python3/dist-packages
ExecStart=/opt/gecko-controller/venv/bin/python3 -m gecko_controller.controller
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target