# debian/postrm
#!/bin/sh
set -e

VENV_DIR="/opt/gecko-controller/venv"

case "$1" in
    remove|purge)
        # Stop and disable the service
        systemctl stop gecko-controller || true
        systemctl disable gecko-controller || true

        # Remove the virtual environment
        rm -rf "$VENV_DIR"
        ;;
esac

#DEBHELPER#
