# debian/postinst
#!/bin/sh
set -e

VENV_DIR="/opt/gecko-controller/venv"
LOG_FILE="/var/log/gecko-controller-install.log"

log() {
    echo "$(date '+%Y-%m-%d %H:%M:%S') - $1" | tee -a "$LOG_FILE"
}

get_python_version() {
    "$VENV_DIR/bin/python3" -c 'import sys; print(f"{sys.version_info.major}.{sys.version_info.minor}")'
}

case "$1" in
    configure)
        log "Starting gecko-controller installation on $(uname -n)"

        # Create virtual environment
        log "Creating virtual environment in $VENV_DIR"
        python3 -m venv "$VENV_DIR"

        # Get Python version for site-packages path
        PYTHON_VERSION=$(get_python_version)
        log "Detected Python version: $PYTHON_VERSION"

        # Upgrade pip and install wheel
        log "Upgrading pip and installing wheel"
        "$VENV_DIR/bin/pip" install --upgrade pip wheel setuptools

        # Install required packages
        log "Installing required packages"
        "$VENV_DIR/bin/pip" install \
            'RPi.GPIO' \
            'smbus2' \
            'Pillow'

        # Verify installations
        log "Verifying installations"
        installed_packages=$("$VENV_DIR/bin/pip" freeze)
        echo "Installed packages:" | tee -a "$LOG_FILE"
        echo "$installed_packages" | tee -a "$LOG_FILE"

        # Enable and start the service
        log "Configuring systemd service"
        systemctl daemon-reload

        log "Enabling gecko-controller service"
        systemctl enable gecko-controller

        log "Starting gecko-controller service"
        systemctl start gecko-controller

        log "Installation completed successfully"
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
esac

#DEBHELPER#